# ansible/roles/amdgpu/templates/etc/profile.d/99-hardware-acceleration.sh.j2
LIBVA_DRIVER_NAME={{ libva_driver_name }}
VDPAU_DRIVER={{ vdpau_driver }}
