# ansible/roles/drivers/templates/drivers/etc/profile.d/99-hardware-acceleration.sh.j2
export LIBVA_DRIVER_NAME={{ libva_driver_name}}
export VDPAU_DRIVER={{ vdpau_driver_name }}
